# untangle-button-led-monitor/debian/untangle-button-led-monitor.postinst
#! /bin/bash

SERVICE_NAME=ut-button-led-monitor
systemctl enable $SERVICE_NAME

exit 0
