# untangle-hardware-linksys-wrt1900acs/debian/untangle-hardware-linksys-wrt1900acs.postrm
#!/bin/bash

SERVICE=untangle-hardware-linksys-wrt1900acs

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    systemctl disable untangle-hardware-linksys-wrt1900acs >/dev/null # lintian...
    ;;
esac

exit 0
