# untangle-hardware-linksys-wrt1900acs/debian/untangle-hardware-asus-ac88u.postrm
#!/bin/bash

SERVICE=untangle-hardware-linksys-wrt1900acs

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    update-rc.d untangle-hardware-linksys-wrt1900acs remove >/dev/null # lintian...
    ;;
esac

exit 0
