# untangle-hardware-linksys-wrt1900acs/debian/untangle-hardware-linksys-wrt1900acs.postinst
#!/bin/sh

systemctl restart untangle-firmware-common

# 12.2.1 upgrade NGFW-10303
# remove file if it exists
rm -f /etc/udev/rules.d/70-persistent-net.rules

# 13.0 remove wan balancer limitation
rm -f /etc/untangle/hardware/firmware.d/50hide

# disable rsyslog by default
systemctl disable rsyslog

rm -f /etc/init.d/${DPKG_MAINTSCRIPT_PACKAGE}

exit 0
