# untangle-hardware-portwell-u150/debian/untangle-hardware-portwell-u150.postrm
#!/bin/bash

SERVICE=untangle-hardware-portwell-u150

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    systemctl disable untangle-hardware-portwell-u150 >/dev/null # lintian...
    ;;
esac

exit 0
