# untangle-ldap-server/debian/untangle-ldap-server.postrm
#!/bin/bash

SERVICE=untangle-slapd

case $1 in
  purge)
    rm -fr /etc/$SERVICE/slapd.conf
    update-rc.d $SERVICE remove >/dev/null
    rmdir --ignore-fail-on-non-empty /etc/$SERVICE
    ;;
esac

exit 0
