# untangle-firmware-common/files/etc/cron.d/sqlite-clean
*/5 * * * * root /usr/bin/sqlite-cleandb.sh | logger -t uvmreports
