# untangle-serial-config/files/etc/systemd/system/serial-getty@ttyS0.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/sbin/agetty --autologin root -s %I 115200,38400,9600 vt102
