# untangle-clamav-config/debian/untangle-clamav-config.cron.daily
#!/bin/bash

nohup /usr/share/untangle/bin/untangle-clamav-config >/dev/null 2>&1 &
exit 0
