# untangle-linux-config/debian/postinst
#!/bin/bash

if [ "$1" = "configure" ] ; then
  # oldVersion will be non-empty if untangle-linux-config has been
  # configured before
  oldVersion="$2" 
fi

ourInit() {
if [ -x /usr/sbin/invoke-rc.d ] && [ ! "`readlink /usr/sbin/invoke-rc.d`" = /bin/true ] ; then
  invoke-rc.d $1 $2
else
  /etc/init.d/$1 $2
fi
}

# attempt to fix bug #7974
grep'HWCLOCKPARM' /etc/default/rcS > /dev/null 2>&1
if [ ! $? = 0 ] ; then 
    echo "HWCLOCKPARS=--directisa" >> /etc/default/rcS
fi
 
# Clean up stuff left around by old version
rm -f /boot/vmlinuz-ut
rm -f /boot/vmlinuz
rm -f /boot/initrd.img

root=
if [ -f /etc/mtab ]; then
    root=`awk '/[ \t]+\/[ \t]+/ {print $1}' /etc/mtab`
fi
if [ -z "$root" -a -f /etc/fstab ]; then
    # Special handling for install-time when /etc/mtab isn't valid
    root=`awk '/[ \t]+\/[ \t]+/ {print $1}' /etc/fstab`
fi

# If we couldn't find a root, don't bother with the rest as we must
# be in some fake install time.
if [ -n "$root" ]; then

    case "$root" in /dev/*)
            rootuuid=`/lib/udev/vol_id -u $root`
            if [ -z "$rootuuid" ]; then
	        rootuuid=`uuidgen`
                # FIXME: Only works for ext2/ext3
                echo "No root UUID found, making one"
	        tune2fs -U "$rootuuid" $root
            fi
            rootdev="UUID=$rootuuid"
            ;;
        *)
            # Fallback
            rootdev="$root"
            ;;
    esac

    # depmod -a ${flavor}
    # Update the initrd for all our kernels, since we changed
    # bootsplash in 6.1
    ls /boot/vmlinuz-2.6*untangle* | while read kernel ; do
      update-initramfs -t -u -k $(echo $kernel | perl -pe 's/.*?-//')
    done

    # Re-create menu.lst
    echo "Installing Untangle /boot/grub/menu.lst."
    currentDefault=$(head -1 /boot/grub/default | perl -pe 's/^(\d).*/$1/')
    [ -z "$currentDefault" ] && currentDefault=0
    currentKernelParams="$(grep -E '^kernel' /boot/grub/menu.lst | head -$((1 + $currentDefault)) | tail -1)"
    case "$currentKernelParams" in
      *ut-video*|"") newDefault=0 ;;
      *noapic*) newDefault=2 ;;
      *ut-restore*) newDefault=3 ;;
      *) newDefault=1 ;;
    esac
    [ -f /boot/grub/menu.lst-pre-untangle ] || cp /boot/grub/menu.lst /boot/grub/menu.lst-pre-untangle
        mkdir -p /boot/grub
    sed -e "s/@ROOTDEV@/$rootdev/" /usr/share/untangle-linux-config/menu.lst > /boot/grub/menu.lst
    /usr/sbin/update-grub
    /usr/sbin/grub-set-default $newDefault
fi

exit 0
