# untangle-linux-config/debian/postinst
#!/bin/bash

set -x

# we don't deal with grub on arm* at all
uname -m | grep -qE '^arm' && exit 0

if [ "$1" = "configure" ] ; then
  # oldVersion will be non-empty if untangle-linux-config has been
  # configured before
  oldVersion="$2" 
fi

grabRootDevice() {
  awk '/^[^#].*[ \t]+\/[ \t]+(type )?ext/ {print $1}' $1
}
 
if [ -f /etc/mtab ]; then
  # use mtab if present
  root=$(grabRootDevice /etc/mtab)
else # no mtab in FAI
  root=$(mount | grabRootDevice)
fi

# If we couldn't find a root, don't bother with the rest as we must
# be in some fake install time.
[ -z "$root" ] && exit 0

case "$root" in
  /dev/*) # let's get a UUID for this device
    rootdev=$(blkid $root | awk '{gsub(/"/, "", $2) ; print $2}') ;;
  *) # fallback: we already have a UUID
    rootdev="$root" ;;
esac

# rewrite fstab
echo "Forcing / in fstab to map to current root device ($rootdev)"
perl -i -pe 's|^[^\s#]+(\s+/\s+)|'$rootdev'$1|' /etc/fstab

currentDefault=$(grub-editenv list | awk -F= '/^saved_entry/ {print $2}')
rm -f /etc/grub.d/10_linux # we use our own 50_custom instead

# move the uefi 'System setup' menu entry to the end of the menu
mv /etc/grub.d/30_uefi-firmware /etc/grub.d/60_uefi-firmware

cp -f /usr/share/untangle-linux-config/grub-default /etc/default/grub
if dpkg -l 'grub-*' 2> /dev/null | grep -qP '^ii .*(grub-cloud-amd64|grub-efi-arm64)' ; then
  mkdir -p /etc/default/grub.d
  cp -f /usr/share/untangle-linux-config/grub-default-extra-cloud /etc/default/grub.d/00_untangle_cloud_serial_console.cfg
fi

get_cpuinfo_field () {
  line=$(grep -m1 ^"$1" /proc/cpuinfo)
  echo "${line##*: }"
}

bay_trail () {
  family=$(get_cpuinfo_field "cpu family")
  model=$(get_cpuinfo_field "model")
  stepping=$(get_cpuinfo_field "stepping")

  if [ $family -eq 6 ] && \
     [ $model -eq 55 ] && \
     [ $stepping -eq 9 ] ; then

    return 0
  else
    return 1
  fi
}

# if this is a bay trail cpu like in the lanner u25, add 'intel_idle.max_cstate=1' to the kernel command line
if bay_trail ; then
  sed -e 's/\(^GRUB_CMDLINE_LINUX=.*\)\("\)/\1 intel_idle.max_cstate=1\2/g' -i /etc/default/grub

  # schedule a reboot if we aren't already running with intel_idle.max_cstate=1
  if ! grep -q intel_idle.max_cstate=1 /proc/cmdline ; then
    echo "Rebooting after apt-get completes..."
    nohup /usr/bin/uvm-restart reboot >> /var/log/uvm/restart.log 2>&1 &
  fi
fi

if [ -z "$currentDefault" ] || [ "$currentDefault" = "default" ] ; then
  echo "Setting default grub option: 0"
  grub-set-default 0
fi

# FIXME: decide what kernel to use upon upgrading from wheezy to jessie
# Reset default (not necessary now, necessary when upgrading kernel)
# currentKernelParams="$(grep -E '^kernel' $menuLst | head -$((1 + $currentDefault)) | tail -1)"
# case "$currentKernelParams" in
#   *force-video-safe*|"") newDefault=1 ;;
#   *noapic*) newDefault=3 ;;
#   *) newDefault=0 ;;
# esac
# grub-set-default $newDefault

# install new grub on update; can't re-use the above UUID=... as
# grub-install doesn't understand that
echo "Configuring grub2..."
if [ -f /boot/grub/device.map ] ; then
  echo "  * Checking in /boot/grub/device.map"
  hd0=$(awk '/hd0/ { print $2 }' /boot/grub/device.map)
  echo "  * found hd0 on '${hd0}'"
  hd0Target=$(readlink -f ${hd0})
  echo "  * target '${hd0Target}'"
fi
if [ ! -e "$hd0Target" ] ; then
  echo "  * No target found in /boot/grub/device.map. Using /"
  if [ -f /etc/mtab ] ; then
    root=$(grabRootDevice /etc/mtab)
  else
    root=$(mount | grabRootDevice)
  fi
  rootTarget=$(readlink -f ${root})
  echo "  * / is on '${root}', which links to '${rootTarget}'"
  hd0Target=$(echo $rootTarget | perl -pe 's/p?\d$//')
  echo "  * target '${hd0Target}'"
fi

echo "Installing grub2... \"${hd0Target}\""

if [ ! -z "${hd0Target}" ] ; then
    grub-install $hd0Target
fi

# update config
echo "Updating grub2..."
update-grub

if [ ! -f /etc/ssl/certs/java/cacerts ] ; then
    echo "Creating default /etc/ssl/certs/java/cacerts file..."
    cp -f /usr/share/untangle-linux-config/default-java-cacerts /etc/ssl/certs/java/cacerts
fi

# enable serial console
if [ `tty` = "/dev/ttyS0" ]; then
  deb-systemd-helper enable serial-getty@ttyS0.service
  deb-systemd-invoke start serial-getty@ttyS0.service
else 
  deb-systemd-invoke stop serial-getty@ttyS0.service
  deb-systemd-helper disable serial-getty@ttyS0.service
  deb-systemd-helper mask serial-getty@ttyS0.service
fi

# enable textui on serial console for zsh or bash shell when root logs in
if ! grep -q ttyS0 /root/.zshrc ; then
    echo 'if [ `tty` = "/dev/ttyS0" ] ; then read ; /usr/share/untangle/bin/ut-textui.py ; fi' >> /root/.zshrc
else
    sed -e 's|.*ttyS0.*|if [ `tty` = "/dev/ttyS0" ] ; then read ; /usr/share/untangle/bin/ut-textui.py ; fi|' -i /root/.zshrc
fi
if ! grep -q ttyS0 /root/.bashrc ; then
    echo 'if [ `tty` = "/dev/ttyS0" ] ; then read ; /usr/share/untangle/bin/ut-textui.py ; fi' >> /root/.bashrc
else
    sed -e 's|.*ttyS0.*|if [ `tty` = "/dev/ttyS0" ] ; then read ; /usr/share/untangle/bin/ut-textui.py ; fi|' -i /root/.bashrc
fi

# 14.0.0 upgrade fix
# overwrite shitty miniupnp scripts
echo -e '#!/bin/sh\nexit 0\n' > /etc/miniupnpd/iptables_init.sh
echo -e '#!/bin/sh\nexit 0\n' > /etc/miniupnpd/ip6tables_init.sh

# enable the systemd service that makes sure we call depmod -a
# before we reboot
deb-systemd-helper enable untangle-linux-config.service

exit 0
