# untangle-hardware-config/debian/untangle-hardware-config.postrm
#!/bin/bash

SERVICE=untangle-hardware-config

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    systemctl disable untangle-hardware-config >/dev/null # lintian...
    ;;
esac

exit 0
