# untangle-spamassassin-update/files/etc/cron.daily/untangle-spamassassin-update
#! /bin/bash

# attempt to fetch the GPG keys anyway
. /usr/share/untangle-spamassassin-update/sa-fetch-keys.sh

# SARE fingerprint
SARE_FINGERPRINT=`gpg ${SARE_GPG_KEY} | awk '/pub/ {gsub(/.*\//, "", $2) ; print $2}'`

if [ -x /usr/bin/sa-update ]; then
    /usr/bin/sa-update --channelfile /usr/share/untangle-spamassassin-update/sa-update-channels.txt --gpgkey $SARE_FINGERPRINT
fi

# restart spamassassin
/etc/init.d/spamassassin restart &> /dev/null &

# sa-update exits with a status to tell us if an update happened.
# don't pass this up to cron.
exit 0
