# untangle-spamassassin-update/files/etc/cron.daily/untangle-spamassassin-update
#! /bin/bash

if [ -x /usr/bin/sa-update ]; then
    /usr/bin/sa-update --channelfile /usr/share/untangle-spamassassin-update/sa-update-channels.txt
fi

# sa-update exits with a status to tell us if an update happened.
# don't pass this up to cron.
exit 0
