# untangle-nflogd/debian/untangle-nflogd.service
[Unit]
Description=Untangle NFlogD

[Service]
Type=simple
ExecStart=/usr/bin/nflogd
Restart=always

[Install]
WantedBy=multi-user.target
