# untangle-firmware-common/files/usr/share/untangle-firmware-common/getty@ttyS0.service
[Service]
Type=simple
ExecStart=
ExecStart=-/sbin/getty -L ttyS0 115200 vt100
