# untangle-pyconnector/debian/untangle-pyconnector.postinst
#!/bin/sh

deb-systemd-helper disable untangle-pyconnector.service

# 14.0 remove old log files
rm -f /var/log/uvm/pyconnector.log*

exit 0
