# untangle-pyconnector/debian/untangle-pyconnector.postinst
#!/bin/sh

systemctl disable untangle-pyconnector

exit 0
