# untangle-pyconnector/debian/untangle-pyconnector.postinst
#!/bin/sh

systemctl disable untangle-pyconnector

# 14.0 remove old log files
rm -f /var/log/uvm/pyconnector.log*

exit 0
