# untangle-restore-tools/debian/untangle-restore-tools.postrm
#! /bin/bash

if [ "$1" = "purge" ] ; then
    rm -f /etc/init.d/untangle-restore-tools
    update-rc.d untangle-restore-tools remove > /dev/null
fi
