# untangle-firmware-common/files/etc/cron.hourly/sqlite
#!/bin/sh

/usr/bin/sqlite-cleandb.sh >> /var/log/uvm/reports.log 2>&1

