# untangle-hardware-lanner-u25w/debian/untangle-hardware-lanner-u25w.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-hardware-lanner-u25w" ]; then
  update-rc.d untangle-hardware-lanner-u25w defaults 10 >/dev/null # early on
  systemctl restart untangle-hardware-lanner-u25w
fi
    
exit 0
