# untangle-hardware-lanner-u25w/debian/untangle-hardware-lanner-u25w.postrm
#!/bin/bash

SERVICE=untangle-hardware-lanner-u25w

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    systemctl disable untangle-hardware-lanner-u25w >/dev/null # lintian...
    ;;
esac

exit 0
