# untangle-tangent-kernel/files/etc/network/if-pre-up.d/untangle-tangent
#! /bin/sh

#modprobe FOO

#/usr/bin/tangent-binary arg1 arg2

exit 0