# untangle-stunnel-config/debian/untangle-stunnel-config.postinst
#! /bin/bash

# enable xinetd for auto start
update-rc.d xinetd enable

exit 0
