# untangle-hardware-portwell-u500/debian/untangle-hardware-portwell-u500.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-hardware-portwell-u500" ]; then
  update-rc.d untangle-hardware-portwell-u500 defaults 10 >/dev/null # early on
  systemctl restart untangle-hardware-portwell-u500
fi
    
exit 0
