# untangle-hardware-lanner-u50/debian/untangle-hardware-lanner-u50.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-hardware-lanner-u50" ]; then
  update-rc.d untangle-hardware-lanner-u50 defaults 10 >/dev/null # early on
  systemctl restart untangle-hardware-lanner-u50
fi
    
exit 0
