# untangle-hardware-lanner-u50/debian/untangle-hardware-lanner-u50.postrm
#!/bin/bash

SERVICE=untangle-hardware-lanner-u50

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    systemctl disable untangle-hardware-lanner-u50 >/dev/null # lintian...
    ;;
esac

exit 0
