# untangle-geoip-database/files/etc/cron.daily/geoip-update
#!/bin/sh
#
# This file is part of the untangle-geoip-database package.  We hook into
# the monthly cron process to download a new database file for the
# Uvm GeographyManager.  We put the new file in the expected location
# with the .update extension and Uvm will take care of the rest.
#

# You can uncomment the following to debug script execution
#set -x

# This is the URL where we download the database update tarball
NETURL="https://downloads.edge.arista.com/download.php"

# This is the location for the temporary working copy we download
TEMPFILE="/tmp/geoip-download.tar.gz"

# This is the name of the update file the Uvm will watch for
UPFILE="/var/cache/untangle-geoip/GeoCountry.update"

# This is the path where we extract the update tarball
TARBALLPATH="/tmp/geoip-update"

# Get the UID so we can pass it in the download request
UID="$(cat /usr/share/untangle/conf/uid)"

MAX_RETRIES=5
SLEEP_SECS=10

cleanup() {
    rm -f "$TEMPFILE"
    rm -rf "$TARBALLPATH"
}

cleanup

# Download with manual retry
i=1
while [ $i -le $MAX_RETRIES ]; do
    echo "Downloading GeoIP DB (attempt $i/$MAX_RETRIES)..."

    if /usr/bin/wget \
        --timeout=30 \
        --read-timeout=30 \
        --retry-connrefused \
        -O "$TEMPFILE" \
        "$NETURL?resource=geoipCountry&uid=$UID"
    then
        break
    fi

    if [ $i -eq $MAX_RETRIES ]; then
        echo "ERROR: Unable to download GeoIP database after retries"
        exit 1
    fi

    i=$((i + 1))
    sleep $SLEEP_SECS
done

# Make sure we actually received the file
if [ ! -s "$TEMPFILE" ]; then
    echo "ERROR: Unable to download GeoIP database file"
    exit 1
fi

if ! gzip -t "$TEMPFILE" >/dev/null 2>&1; then
    echo "ERROR: Downloaded file is not a valid gzip archive"
    exit 1
fi

# Extract
mkdir -p "$TARBALLPATH" || {
    echo "ERROR: Unable to create extraction directory"
    exit 1
}

if ! tar xzf "$TEMPFILE" -C "$TARBALLPATH" --strip=1; then
    echo "ERROR: Unable to extract GeoIP update file"
    cleanup
    exit 1
fi

# This is the full path to the database file extracted from the update tarball
TEMPDBFILE="$(find "$TARBALLPATH" -name 'Geo*Country.mmdb' -type f | head -1)"

if [ -z "$TEMPDBFILE" ] || [ ! -s "$TEMPDBFILE" ]; then
    echo "ERROR: GeoIP database file not found or empty after extraction"
    cleanup
    exit 1
fi

# Publish update
rm -f "$UPFILE" || {
    echo "ERROR: Unable to remove existing GeoIP update file"
    cleanup
    exit 1
}

cp "$TEMPDBFILE" "$UPFILE" || {
    echo "ERROR: Unable to publish GeoIP update"
    cleanup
    exit 1
}

cleanup
exit 0
