# untangle-kiosk/files/usr/share/untangle-kiosk/getty@tty1.service
[Service]
Type=simple
ExecStart=
ExecStart=-/sbin/agetty --autologin kiosk --noclear %I 38400 linux
