# untangle-hardware-lanner-u25/debian/untangle-hardware-lanner-u25.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-hardware-lanner-u25" ]; then
  update-rc.d untangle-hardware-lanner-u25 defaults 10 >/dev/null # early on
  systemctl restart untangle-hardware-lanner-u25
fi
    
exit 0
