# untangle-hardware-asus-ac88u/debian/untangle-hardware-asus-ac88u.postrm
#!/bin/bash

SERVICE=untangle-hardware-asus-ac88u

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    update-rc.d untangle-hardware-asus-ac88u remove >/dev/null # lintian...
    ;;
esac

exit 0
