# untangle-hardware-asus-ac88u/debian/untangle-hardware-asus-ac88u.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-firmware-common" ]; then
  systemctl restart untangle-firmware-common
fi

# disable rsyslog by default on linksys
systemctl disable rsyslog

# remove old 12.1 conf
rm -rf /etc/untangle/hardware/asus

# remove old 12.1 startup file
systemctl disable untangle-hardware-asus-ac88u
rm -f /etc/init.d/untangle-hardware-asus-ac88u

exit 0
