# untangle-hardware-jetway-u50w/debian/untangle-hardware-jetway-u50w.postrm
#!/bin/bash

SERVICE=untangle-hardware-jetway-u50w

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    systemctl disable untangle-hardware-jetway-u50w >/dev/null # lintian...
    ;;
esac

exit 0
