# untangle-hardware-jetway-u50w/debian/untangle-hardware-portwell-u5.postrm
#!/bin/bash

SERVICE=untangle-hardware-jetway-u50w

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    update-rc.d untangle-hardware-jetway-u50w remove >/dev/null # lintian...
    ;;
esac

exit 0
