# untangle-hardware-jetway-u50w/debian/untangle-hardware-jetway-u50w.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-hardware-jetway-u50w" ]; then
  update-rc.d untangle-hardware-jetway-u50w defaults 10 >/dev/null # early on
  ourInit untangle-hardware-jetway-u50w restart
fi
    
exit 0
