# untangle-broadcom-wireless/hier/etc/network/if-up.d/broadcom-wireless
#! /bin/bash

brcm-wifi start eth1
brcm-wifi start eth2

exit 0
