# untangle-broadcom-wireless/hier/etc/network/if-up.d/broadcom-wireless
#! /bin/bash

brcm-wifi.sh start eth1
brcm-wifi.sh start eth2

exit 0

