# untangle-softflowd-config/debian/untangle-softflowd-config.postinst
#! /bin/bash

DEFAULT_FILE="/etc/default/softflowd"

systemctl disable softflowd

# Only for first time installation.
if [ -z "$2" ] ; then
    systemctl stop softflowd
fi

exit 0
