# untangle-softflowd-config/debian/untangle-softflowd-config.postinst
#! /bin/bash

DEFAULT_FILE="/etc/default/softflowd"

deb-systemd-helper disable softflowd.service

# Only for first time installation.
if [ -z "$2" ] ; then
    deb-systemd-invoke stop softflowd
fi

exit 0
