# untangle-gateway/debian/postinst
#! /bin/sh

# do not start sshd by default if this is the first install
if [ "$1" = "upgrade" ] ; then
  update-rc.d -f ssh remove
  /etc/init.d/ssh stop
fi

exit 0
