# untangle-net-alpaca/debian/untangle-net-alpaca-iptables.postrm
#!/bin/bash

SERVICE=untangle-net-alpaca-iptables

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    update-rc.d $SERVICE remove >/dev/null
    ;;
esac

exit 0
