# untangle-hardware-lanner-u500/debian/untangle-hardware-lanner-u500.init
#! /bin/sh

### BEGIN INIT INFO
# Provides:          untangle-hardware-lanner-u500
# Required-Start:    $local_fs $network
# Required-Stop:
# Default-Start:     S
# Default-Stop:      
### END INIT INFO

# Author: David Han <dhan@untangle.com>

. /lib/lsb/init-functions

DESC="Untangle Hardware Lanner u500"
PARTS_DIR="/etc/untangle/hardware/lanner/u500.d"
LOG_FILE="/var/log/uvm/lanner-u500.log"

start() {
  log_daemon_msg "Starting $DESC"
  log_progress_msg "untangle-hardware-lanner-u500"    
  run-parts --lsbsysinit -v $PARTS_DIR >> $LOG_FILE 2>&1
  RET=$?
  log_end_msg $RET
  [ $RET = 0 ] || exit $RET
}

stop() {
  log_daemon_msg "Stopping $DESC"
  log_progress_msg "untangle-hardware-lanner-u500"    
  log_end_msg $RET
}

case "$1" in
    start)
        start
        ;;
    stop)
        stop
        ;;
    reload|restart|force-reload)
        stop
        start
        ;;
    *)
        /bin/echo "Usage: $0 {start|stop|restart|force-reload}" >&2
        exit 1
        ;;
esac

exit 0
