# untangle-google-safebrowsing/files/etc/cron.d/untangle-google-safebrowsing
* * */1 * * root /usr/share/untangle-google-safebrowsing/bin/update_malware_list.sh > /dev/null 2>&1
* * */1 * * root /usr/share/untangle-google-safebrowsing/bin/update_phish_list.sh > /dev/null 2>&1
