# untangle-hardware-linksys-wrt3200acm/debian/untangle-hardware-linksys-wrt3200acm.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-firmware-common" ]; then
  systemctl restart untangle-firmware-common
fi

# 12.2.1 upgrade NGFW-10303
# remove file if it exists
rm -f /etc/udev/rules.d/70-persistent-net.rules

# 13.0 remove wan balancer limitation
rm -f /etc/untangle/hardware/firmware.d/50hide

# disable rsyslog by default on linksys
systemctl disable rsyslog

exit 0
