# untangle-softflowd-config/debian/untangle-softflowd-config/DEBIAN/postinst
#! /bin/bash

DEFAULT_FILE="/etc/default/softflowd"

update-rc.d -f softflowd disable

# Only for first time installation.
if [ -z "$2" ] ; then
    service softflowd stop
fi

exit 0
