# untangle-net-alpaca/files/etc/network/if-up.d/alpaca-update-address
#!/bin/dash

[ "$IFACE" != "update" ] && {
    exit 0
}

run-parts /etc/untangle-net-alpaca/scripts/update-address.d

exit 0


