# untangle-hardware-lanner-m3000/debian/untangle-hardware-lanner-m3000.postrm
#!/bin/bash

SERVICE=untangle-hardware-lanner-m3000

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    systemctl disable untangle-hardware-lanner-m3000 >/dev/null # lintian...
    ;;
esac

exit 0
