# untangle-hardware-lanner-m3000/debian/untangle-hardware-lanner-m3000.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-hardware-lanner-m3000" ]; then
  update-rc.d untangle-hardware-lanner-m3000 defaults 10 >/dev/null # early on
  systemctl restart untangle-hardware-lanner-m3000
fi
    
exit 0
