# untangle-net-alpaca/files/etc/cron.d/untangle-net-alpaca
*/5 * * * * root  /etc/untangle-net-alpaca/untangle-qos-rrdtool update
*/10 * * * * root  /etc/untangle-net-alpaca/untangle-qos-rrdtool graph > /dev/null
