# untangle-hardware-portwell-m1500/debian/untangle-hardware-portwell-m1500.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-hardware-portwell-m1500" ]; then
  update-rc.d untangle-hardware-portwell-m1500 defaults 10 >/dev/null # early on
  ourInit untangle-hardware-portwell-m1500 restart
fi
    
exit 0
