# untangle-hardware-portwell-m1500/debian/untangle-hardware-portwell-m1500.postinst
#!/bin/sh

ourInit() {
if [ -x /usr/sbin/invoke-rc.d ] && [ ! "`readlink /usr/sbin/invoke-rc.d`" = /bin/true ] ; then
  invoke-rc.d $1 $2
else
  /etc/init.d/$1 $2
fi
}

if [ -x "/etc/init.d/untangle-hardware-portwell-m1500" ]; then
  update-rc.d untangle-hardware-portwell-m1500 defaults 10 >/dev/null # early on
  ourInit untangle-hardware-portwell-m1500 restart
fi
    
exit 0
