# untangle-hardware-qotom-u25xw/debian/postrm
#!/bin/bash

SERVICE=${DPKG_MAINTSCRIPT_PACKAGE}

case $1 in
  purge)
    rm -f /etc/init.d/$SERVICE
    update-rc.d $SERVICE remove >/dev/null
    ;;
esac

exit 0
