# untangle-hardware-turris-omnia/debian/untangle-hardware-turris-omnia.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-firmware-common" ]; then
  systemctl restart untangle-firmware-common
fi

# disable rsyslog by default
systemctl disable rsyslog

exit 0
