# untangle-hardware-turris-omnia/debian/untangle-hardware-turris-omnia.postinst
#! /bin/bash

deb-systemd-invoke restart untangle-firmware-common

# disable rsyslog by default
deb-systemd-helper disable rsyslog.service

rm -f /etc/init.d/${DPKG_MAINTSCRIPT_PACKAGE}

exit 0
