# restore-mv/debian/restore-mv.postrm
#!/bin/sh

if [ "$1" = "purge" ] ; then
    rm -f /etc/init.d/restore-mv
    update-rc.d restore-mv remove > /dev/null
fi
