# untangle-hardware-portwell-u10/debian/untangle-hardware-portwell-u10.postinst
#!/bin/sh

if [ -x "/etc/init.d/untangle-hardware-portwell-u10" ]; then
  update-rc.d untangle-hardware-portwell-u10 defaults 10 >/dev/null # early on
  systemctl restart untangle-hardware-portwell-u10
fi

exit 0
