# untangle-snort-config/debian/untangle-snort-config.files
/etc/init.d/*
/etc/logrotate.d/*
/etc/snort/*
/etc/rsyslog.d/*

